# scripts/linux/build-it-agent.service
[Unit]
Description=Build It Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/build-it-agent
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
